# scripts/60-can.rules
ACTION=="add", ATTR{idVendor}=="0c72", RUN+="/bin/reset_can.sh"
